# systemd/rkhunter.timer
[Unit]
description=Daily rootkit scan

[Timer]
OnCalendar=00:00:00
Persistent=true

[Install]
WantedBy=timers.target
